# deploy/systemd/avapharmacy-production.service
[Unit]
Description=AvaPharmacy production Gunicorn service
After=network.target

[Service]
User=ava
Group=www-data
WorkingDirectory=/home/ava/production/backend
EnvironmentFile=/home/ava/production/backend/.env
ExecStart=/home/ava/production/backend/venv/bin/gunicorn avapharmacy.wsgi:application --bind 127.0.0.1:8002 --workers 3 --timeout 120
Restart=always
RestartSec=5

[Install]
WantedBy=multi-user.target
